# applications/myapp/cron/crontab
#crontab
@reboot  root *applications/myapp/cron/tablebar.py
30 9 * * *  root *applications/myapp/cron/sendpcmlog.py
