# applications/myapp/cron/crontab
#crontab
@reboot  root *applications/myapp/cron/tablebar.py
@reboot  root *applications/myapp/cron/sendpcmlog.py
